# system_files/etc/profile.d/gamerslop.sh
export PROTON_ENABLE_WOW64=1
export PROTON_USE_WOW64=1
export PROTON_ENABLE_NTSYNC=1
export PROTON_USE_NTSYNC=1
export PROTON_ENABLE_WAYLAND=1
export PROTON_USE_WAYLAND=1
export PROTON_PREFER_SDL=1
